# conf/supervisord.conf
; supervisor config file for terane test cluster.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html

[supervisord]
directory = ../run
logfile = ../run/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = debug
childlogdir = ../run
pidfile = ../run/supervisord.pid
nodaemon = false
identifier = terane

[program:cassandra]
directory = ../ext/lib/cassandra
command = bash ./bin/cassandra -f
autostart = false

[program:zookeeper]
directory = ../ext/lib/zookeeper
command = bash ./bin/zkServer.sh start-foreground
autostart = false

[program:terane-test1]
directory = ../
command = java -Dterane.config.file=conf/terane-test1.conf -jar target/scala-2.10/terane-indexer_2.10-0.1-one-jar.jar
autostart = false

[program:terane-test2]
directory = ../
command = java -Dterane.config.file=conf/terane-test2.conf -jar target/scala-2.10/terane-indexer_2.10-0.1-one-jar.jar
autostart = false

[program:terane-test3]
directory = ../
command = java -Dterane.config.file=conf/terane-test3.conf -jar target/scala-2.10/terane-indexer_2.10-0.1-one-jar.jar
autostart = false

[supervisorctl]
serverurl = http://127.0.0.1:9001
prompt = terane
history_file = ~/.sc_history

[inet_http_server]
port = 127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


